# roles/notes_rotation/templates/notes-journal-rotate@.timer.j2
[Unit]
Description=Rotate Notes Journal

[Timer]
OnCalendar=*-*-* 23:50:00

[Install]
WantedBy=default.target

