# roles/gitea_repos_clone/templates/git-fetch-@.timer.j2
[Unit]
Description=Run fetch in git repo
Documentation=man:git

[Timer]
OnCalendar=daily
RandomizedDelaySec=4h
Persistent=false

[Install]
WantedBy=timers.target
